# camillagui.service
[Unit]
Description=CamillaDSP Backend and GUI
After=default.target

[Service]
Type=idle
User=username
WorkingDirectory=~
ExecStart=/home/username/camilladsp/camillagui_backend/camillagui_backend

[Install]
WantedBy=default.target